# rootfs/etc/s6-overlay/s6-rc.d/hcloud-ddns-api/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Hetzner DNS DDNS Addon
# API service - serves the force-update web UI
# ==============================================================================

bashio::log.info "Starting Hetzner DNS DDNS API server..."

exec python3 /usr/bin/hcloud-ddns-api.py
