# rootfs/etc/s6-overlay/s6-rc.d/hcloud-ddns/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Hetzner DNS DDNS Addon
# Main service - runs the DDNS update script
# ==============================================================================

bashio::log.info "Starting Hetzner DNS DDNS service..."

exec /usr/bin/hcloud-ddns.sh
