# rootfs/etc/s6-overlay/s6-rc.d/init-hcloud/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Hetzner Cloud DDNS Addon
# Init script - validates API key and configuration
# ==============================================================================

bashio::log.info "Initializing Hetzner DNS DDNS addon..."

# Validate API key
API_KEY=$(bashio::config 'api_key')
if [ -z "${API_KEY}" ]; then
    bashio::log.fatal "API key is not configured"
    exit 1
fi

# Store API key securely
mkdir -p /root/.config
chmod 700 /root/.config
echo "${API_KEY}" > /root/.config/hetzner-dns-token
chmod 600 /root/.config/hetzner-dns-token

bashio::log.info "API key configured securely"

# Validate zone ID
ZONE_ID=$(bashio::config 'zone_id')
if [ -z "${ZONE_ID}" ]; then
    bashio::log.fatal "Zone ID is not configured"
    exit 1
fi

# Validate domain
DOMAIN=$(bashio::config 'domain')
if [ -z "${DOMAIN}" ]; then
    bashio::log.fatal "Domain is not configured"
    exit 1
fi

# Test API connection and fetch zone name for use by main script
bashio::log.info "Testing Hetzner Cloud API connection..."
RESPONSE=$(curl -s -w "\n%{http_code}" \
    -H "Authorization: Bearer ${API_KEY}" \
    "https://api.hetzner.cloud/v1/zones/${ZONE_ID}")

HTTP_CODE=$(echo "${RESPONSE}" | tail -1)
BODY=$(echo "${RESPONSE}" | head -n -1)

if [ "${HTTP_CODE}" = "200" ]; then
    ZONE_NAME=$(echo "${BODY}" | jq -r '.zone.name // empty')
    if [ -z "${ZONE_NAME}" ]; then
        bashio::log.fatal "Could not read zone name from API response"
        exit 1
    fi
    echo "${ZONE_NAME}" > /root/.config/hetzner-zone-name
    chmod 600 /root/.config/hetzner-zone-name
    bashio::log.info "API connection successful (zone: ${ZONE_NAME})"
elif [ "${HTTP_CODE}" = "401" ] || [ "${HTTP_CODE}" = "403" ]; then
    bashio::log.fatal "Hetzner Cloud API authentication failed (HTTP ${HTTP_CODE}). Check your API key and permissions."
    exit 1
elif [ "${HTTP_CODE}" = "404" ]; then
    bashio::log.fatal "Hetzner DNS zone not found (HTTP 404). Check your zone ID."
    exit 1
else
    bashio::log.warning "Could not verify DNS zone (HTTP ${HTTP_CODE}). This may be a temporary issue, but configuration should be checked."
fi

bashio::log.info "Initialization complete"
