# hcloud-ddns/rootfs/etc/s6-overlay/s6-rc.d/hcloud-ddns/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Hetzner Cloud DDNS Addon
# Finish script for hcloud-ddns service
# ==============================================================================

bashio::log.info "Hetzner Cloud DDNS service stopped"
